# jb-rotary.timer
[Unit]
Desciption=Delay jb-rotary till ALSA available

[Timer]
OnStartupSec=20

[Install]
WantedBy=multi-user.target
